# crontab
* * * * * /app/.venv/bin/python /app/rbcua.py "tmp/data.csv"
